# openhd.service
[Unit]
Description=OpenHD 

[Service]
User=root
ExecStart=/usr/local/bin/OpenHD
Restart=always


[Install]
WantedBy=multi-user.target
