# overlay/etc/systemd/system/osd.service
[Unit]
Description=OSD
After=multi-user.target

[Service]
Type=simple
ExecStart=/usr/local/bin/osd
User=root

[Install]
WantedBy=multi-user.target

